# deploy/oracle/nexora-api.service
[Unit]
Description=NEXORA FastAPI backend
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/srv/nexora/app
EnvironmentFile=/srv/nexora/app/.env.production
ExecStart=/srv/nexora/app/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
